# deploy/com.palate.prod.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.palate.mcp</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd /Volumes/xpg_usb4/prod/palate/current &amp;&amp; exec /Library/Frameworks/Python.framework/Versions/3.14/bin/python3 -m palate.server</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/oric</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Volumes/xpg_usb4/prod/palate/current/.venv/bin:/Library/Frameworks/Python.framework/Versions/3.14/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>PYTHONPATH</key>
    <string>/Volumes/xpg_usb4/prod/palate/current:/Volumes/xpg_usb4/prod/palate/current/.venv/lib/python3.14/site-packages</string>
    <key>PALATE_ENV_PATH</key>
    <string>/Volumes/xpg_usb4/prod/palate/shared/.env</string>
    <key>PALATE_TRANSPORT</key>
    <string>streamable-http</string>
    <key>PALATE_HOST</key>
    <string>127.0.0.1</string>
    <key>PALATE_PORT</key>
    <string>8787</string>
    <key>PALATE_MCP_PATH</key>
    <string>/mcp</string>
    <key>PALATE_DB_PATH</key>
    <string>/Volumes/xpg_usb4/prod/palate/shared/data/palate.sqlite</string>
    <key>PALATE_MODEL</key>
    <string>gpt-5.4-mini</string>
    <key>PALATE_BACKUP_ENABLED</key>
    <string>1</string>
    <key>PALATE_BACKUP_DIR</key>
    <string>/Volumes/xpg_usb4/prod/palate/shared/backups</string>
    <key>PALATE_BACKUP_RETENTION_DAYS</key>
    <string>31</string>
    <key>PALATE_BACKUP_INTERVAL_SECONDS</key>
    <string>86400</string>
    <key>PALATE_BACKUP_GOOGLE_DRIVE_ENABLED</key>
    <string>1</string>
    <key>PALATE_GOOGLE_DRIVE_FOLDER_NAME</key>
    <string>backup/palate</string>
    <key>PALATE_GOOGLE_CREDENTIALS_PATH</key>
    <string>/Volumes/xpg_usb4/prod/palate/shared/secrets/google-oauth-client.json</string>
    <key>PALATE_GOOGLE_TOKEN_PATH</key>
    <string>/Volumes/xpg_usb4/prod/palate/shared/secrets/google-token.json</string>
    <key>PALATE_AUTH_ENABLED</key>
    <string>1</string>
    <key>PALATE_PUBLIC_BASE_URL</key>
    <string>https://palate.dceb.net</string>
    <key>PALATE_PUBLIC_MCP_PATH</key>
    <string>/mcp</string>
    <key>PALATE_AUTH_SCOPES</key>
    <string>palate.access</string>
    <key>PALATE_AUTH_PASSWORD_FILE</key>
    <string>/Volumes/xpg_usb4/prod/palate/shared/secrets/palate-auth-password</string>
    <key>PALATE_AUTH_STATE_PATH</key>
    <string>/Volumes/xpg_usb4/prod/palate/shared/secrets/palate-oauth.json</string>
    <key>PALATE_ALLOWED_HOSTS</key>
    <string>127.0.0.1,127.0.0.1:8787,localhost,localhost:8787,modal.tail63a6b7.ts.net,palate.dceb.net</string>
    <key>PALATE_ALLOWED_ORIGINS</key>
    <string>https://modal.tail63a6b7.ts.net,https://palate.dceb.net</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/oric/Library/Logs/palate-prod.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/oric/Library/Logs/palate-prod.err.log</string>
</dict>
</plist>
